# supervisor/apache.conf
[program:apache]
command=/usr/sbin/apache2