# supervisor/apache.conf
[program:apache2]
command=/usr/sbin/apache2 -DFOREGROUND
redirect_stderr=true